# os/etc/service/wsgi/finish
#!/bin/sh 

uwsgi --stop /run/uwsgi.pid
