# os/etc/crontabs/root
# DO NOT EDIT THIS FILE - edit the master and reinstall.
# (/tmp/crontab.nS0S9F/crontab installed on Wed Oct  5 09:29:23 2022)
# (Cron version -- $Id: crontab.c,v 2.13 1994/01/17 03:20:37 vixie Exp $)
0       8       */1     *       *       python /adcm/python/manage.py logrotate --target all
0       10      */1     *       *       python /adcm/python/manage.py clearaudit
*/1     *       *       *       *       python /adcm/python/manage.py run_ldap_sync
0       0       *       *       1       python /adcm/python/manage.py collect_statistics --mode send
