# os/etc/crontabs/root
# Licensed under the Apache License, Version 2.0 (the "License");
# you may not use this file except in compliance with the License.
# You may obtain a copy of the License at
#
#      http://www.apache.org/licenses/LICENSE-2.0
#
# Unless required by applicable law or agreed to in writing, software
# distributed under the License is distributed on an "AS IS" BASIS,
# WITHOUT WARRANTIES OR CONDITIONS OF ANY KIND, either express or implied.
# See the License for the specific language governing permissions and
# limitations under the License.

# do daily/weekly/monthly maintenance
# min   hour    day     month   weekday command
*/15    *       *       *       *       run-parts /etc/periodic/15min
0       *       *       *       *       run-parts /etc/periodic/hourly
0       2       *       *       *       run-parts /etc/periodic/daily
0       3       *       *       6       run-parts /etc/periodic/weekly
0       5       1       *       *       run-parts /etc/periodic/monthly
0       8       */1     *       *       python /adcm/python/manage.py logrotate --target all > /proc/1/fd/1 2>&1
0       10      */1     *       *       python /adcm/python/manage.py clearaudit > /proc/1/fd/1 2>&1
*/1     *       *       *       *       python /adcm/python/manage.py run_ldap_sync > /proc/1/fd/1 2>&1
*/1     *       *       *       *       python /adcm/python/manage.py simplecommand > /proc/1/fd/1 2>&1
